# surlatablo/config/cron.daily/daily_cpFromTablo
#!/bin/sh

# script abstraction to make this container more flexible and
#  minimize the need to rebuild.
/config/scripts/cpFromTablo >> /var/log/cron.log 2>&1

exit 0
