# surlatablo/root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "Starting Cron Daemon."
exec cron -f
